# tenxsom-scheduler.service
[Unit]
Description=TenxsomAI Daily Content Scheduler
After=network.target

[Service]
Type=simple
User=golde
Group=golde
WorkingDirectory=/home/golde/tenxsom-ai-vertex
ExecStart=/home/golde/tenxsom-ai-vertex/venv/bin/python /home/golde/tenxsom-ai-vertex/daily_content_scheduler.py --daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/golde/tenxsom-ai-vertex
Environment=GOOGLE_APPLICATION_CREDENTIALS=/home/golde/tenxsom-ai-vertex/credentials/tenxsom-ai-service-account.json
Environment=CLOUD_TASKS_WORKER_URL=https://tenxsom-video-worker-hpkm6siuqq-uc.a.run.app/process_video_job

[Install]
WantedBy=multi-user.target