# systemd/tenxsom-scheduler.service
[Unit]
Description=TenxsomAI Daily Content Scheduler
Documentation=file:///home/golde/tenxsom-ai-vertex/PRODUCTION-INTEGRATION-SUMMARY.md
After=network-online.target
Wants=network-online.target

[Service]
# Security settings
User=golde
Group=golde
NoNewPrivileges=true
PrivateTmp=true

# Working directory and environment
WorkingDirectory=/home/golde/tenxsom-ai-vertex
EnvironmentFile=/home/golde/tenxsom-ai-vertex/production-config.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/golde/tenxsom-ai-vertex"

# Process management
Type=exec
ExecStartPre=/bin/bash -c 'test -f /home/golde/tenxsom-ai-vertex/production-config.env || (echo "ERROR: production-config.env not found!" && exit 1)'
ExecStart=/usr/bin/python3 /home/golde/tenxsom-ai-vertex/daily_content_scheduler.py --daemon --production
ExecReload=/bin/kill -USR1 $MAINPID

# Restart policy with exponential backoff
Restart=on-failure
RestartSec=30s
RestartPreventExitStatus=0
StartLimitInterval=600s
StartLimitBurst=5

# Resource limits
MemoryLimit=4G
CPUQuota=50%
TasksMax=200

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60s
SendSIGKILL=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenxsom-scheduler

# Health checking
ExecStartPost=/bin/sleep 10
ExecStartPost=/bin/bash -c 'ps aux | grep -v grep | grep "daily_content_scheduler" || exit 1'

[Install]
WantedBy=multi-user.target