# tenxsom_flow_engine/tenxsom-flow-worker.service
[Unit]
Description=Tenxsom AI Headless Flow Worker
After=network.target redis.service

[Service]
User=golde
Group=golde
WorkingDirectory=/home/golde/tenxsom-ai-vertex/tenxsom_flow_engine/
# Use the system Python (production environment)
ExecStart=/usr/bin/python3 worker.py
Restart=on-failure
RestartSec=10

# Logging will be automatically handled by journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenxsom-flow-worker

[Install]
WantedBy=multi-user.target